# .config/systemd/user/stacktile.service
[Unit]
Description=stacktile
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/stacktile \
	--per-tag-config \
	--primary-count 1 \
	--primary-ratio 0.4 \
	--primary-sublayout columns \
	--primary-position left \
	--secondary-count 1 \
	--secondary-ratio 0.7 \
	--secondary-sublayout rows \
	--remainder-sublayout stack
Restart=on-failure

[Install]
WantedBy=graphical-session.target
